# ops/systemd/skyresq-payload-servo.service
[Unit]
Description=SkyResQ - Payload servo (GPIO 16 / pin 36)
PartOf=skyresq-core.target
After=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/Drone/ros2_ws
EnvironmentFile=-/home/pi/Drone/ops/config
ExecStart=/bin/bash -c 'source /home/pi/Drone/ops/env.sh && exec ros2 run drone_vision payload_servo'
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=skyresq-payload

[Install]
WantedBy=skyresq-core.target
